# build/config/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
silent=false
loglevel=info
minfds=1024
minprocs=200

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[inet_http_server]
; Optional supervisord web interface - only enabled if SUPERVISORD_UI_PORT env var is set
; Set SUPERVISORD_UI_PORT to enable (e.g., 9001)
port=%(ENV_SUPERVISORD_UI_PORT)s
username=%(ENV_SUPERVISORD_UI_USER)s
password=%(ENV_SUPERVISORD_UI_PASS)s

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:settings]
## make sure setting for this container are initialized.
command=init_container_settings
priority=5
autorestart=false
autostart=true
exitcodes=0
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:xvfb]
## start virtual frame buffer.
# creates screen screennum and sets its width, height, and depth to W, H, and D respectively.
command=start_xvfb
autorestart=true
autostart=true
priority=10
startsecs=10
startretries=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:x11vnc]
## start VNC server.
# display: X11 server display to connect to.
# forever: Keep listening for more connections rather than exiting as soon as the first client(s) disconnect.
# shared: VNC display is shared, i.e. more than one viewer can connect at the same time.
# bg: Go into the background after screen setup. Messages to stderr are lost unless -o logfile is used.
# noipv6: Do not try to use IPv6 for any listening or connecting sockets.
# logappend: Write stderr messages to file logfile instead of to the terminal.
command=start_vnc
autorestart=true
autostart=true
priority=20
startsecs=10
startretries=3
depends_on=xvfb
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:ibc]
## start IBC (Interactive Brokers Controller).
# This runs the IBC application which controls TWS/Gateway
command=start_ibc
autorestart=true
autostart=true
priority=30
startsecs=30
startretries=3
stopwaitsecs=60
stopsignal=TERM
# kill java child processes and threads.
killasgroup=true
stopasgroup=true
depends_on=xvfb,x11vnc,settings
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

#[program:socat]
## Port forwarding with socat
#command=start_socat
#priority=40
#autorestart=true
#autostart=true
#depends_on=ibc
#stdout_logfile=/dev/stdout
#stdout_logfile_maxbytes=0
#stderr_logfile=/dev/stderr
#stderr_logfile_maxbytes=0

#[program:ssh-tunnel]
## SSH tunnel for secure remote access
#command=start_ssh
#priority=50
#autorestart=false
#autostart=true
#exitcodes=0
#startsecs=0
#depends_on=ibc
#stdout_logfile=/dev/stdout
#stdout_logfile_maxbytes=0
#stderr_logfile=/dev/stderr
#stderr_logfile_maxbytes=0
