# build/config/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
silent=false
loglevel=info

[unix_http_server]
file=/tmp/supervisor.sock

[inet_http_server]
; Optional supervisord web interface - only enabled if HTTP_SERVER_PORT env var is set
; Set HTTP_SERVER_PORT to enable (e.g., 9001)
port=127.0.0.1:%(ENV_HTTP_SERVER_PORT:-9001)s
username=%(ENV_HTTP_SERVER_USER:-)s
password=%(ENV_HTTP_SERVER_PASS:-)s

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:settings]
## make sure setting for this container are initialized.
command=init_container_settings
priority=5
autorestart=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:xvfb]
## start virtual frame buffer.
# creates screen screennum and sets its width, height, and depth to W, H, and D respectively.
command=start_xvfb
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:x11vnc]
## start VNC server.
# display: X11 server display to connect to.
# forever: Keep listening for more connections rather than exiting as soon as the first client(s) disconnect.
# shared: VNC display is shared, i.e. more than one viewer can connect at the same time.
# bg: Go into the background after screen setup. Messages to stderr are lost unless -o logfile is used.
# noipv6: Do not try to use IPv6 for any listening or connecting sockets.
# logappend: Write stderr messages to file logfile instead of to the terminal.
command=start_vnc
autorestart=true
priority=20
depends_on=xvfb
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:ibc]
## start IBC (Interactive Brokers Controller).
# This runs the IBC application which controls TWS/Gateway
command=start_ibc
autorestart=true
priority=30
startsecs=30
stopwaitsecs=30
depends_on=xvfb,settings
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:socat]
## Port forwarding with socat
command=start_socat
priority=40
autorestart=true
autostart=true
depends_on=ibc
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:ssh-tunnel]
## SSH tunnel for secure remote access
command=start_ssh
priority=50
autorestart=true
autostart=true
depends_on=ibc
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
